# config/proximity.conf
# Upstart configuration

description "proximity"
author "Brad Berger"

start on (local-filesystems and net-device-up)
stop on shutdown

respawn

exec sudo /usr/local/bin/proximity